# pi-gen/stage3/02-setup-usb/files/99-usb-umount.rules
SUBSYSTEM=="block", KERNEL=="sd?1", ACTION=="remove", RUN+="/usr/local/bin/usb-umount.sh %k"
