# pi-gen/stage3/02-setup-usb/files/99-usb-mount.rules
SUBSYSTEM=="block", KERNEL=="sd?1", ACTION=="add", ENV{ID_FS_TYPE}=="vfat", RUN+="/usr/local/bin/usb-mount.sh %k"
